# starlink-exporter.service
[Unit]
Description=Starlink Exporter - Prometheus exporter for Starlink metrics
Documentation=https://github.com/brendanbank/starlink_exporter
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/starlink_exporter
Restart=on-failure
RestartSec=10s
User=root
StandardOutput=journal
StandardError=journal
SyslogIdentifier=starlink-exporter

[Install]
WantedBy=multi-user.target
