# debian/starlink-exporter.prerm
#!/bin/sh
set -e

case "$1" in
    remove|deconfigure)
        # Stop the service
        systemctl stop starlink-exporter.service || true
        systemctl disable starlink-exporter.service || true
        ;;
esac

#DEBHELPER#

exit 0
